# bsd/openbsd/configs/rc.d/qbtnox
#!/bin/ksh

daemon="/usr/local/bin/qbittorrent-nox -d"

. /etc/rc.d/rc.subr

rc_cmd $1
