# linux/scripts/tty-colors-light.service
[Unit]
Description=TTY light colors

[Service]
Type=oneshot
ExecStart=/usr/bin/tty-colors-light

[Install]
WantedBy=getty.target
