# linux/scripts/tty-colors-dark.service
[Unit]
Description=TTY dark colors

[Service]
Type=oneshot
ExecStart=/usr/bin/tty-colors-dark

[Install]
WantedBy=getty.target
