# linux/scripts/tty-colors.service
[Unit]
Description=Change tty colors

[Service]
Type=oneshot
ExecStart=/usr/bin/tty-colors

[Install]
WantedBy=getty.target
