# scripts/tty-colors.service
[Unit]
Description=TTY Colors

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/tty-colors

[Install]
WantedBy=multi-user.target
